# root/etc/services.d/wallet/run
#!/usr/bin/with-contenv bash

if [ "$WALLET" == "true" ] && ! pgrep -x "chia_wallet" >/dev/null; then
    cd /chia
    . ./activate
    echo "*** STARTING chia_wallet ***"
    chia start wallet
fi